# competition-test/deploy/rsi2_qqq.service
[Unit]
Description=RSI-2 QQQ Mean Reversion (5%)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Algo-Trade-UMN-Team4
ExecStart=/home/pi/Algo-Trade-UMN-Team4/.venv/bin/python rsi2_qqq.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
