# crypto/deploy/dashboard.service
[Unit]
Description=Crypto Trader Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/algo-trade/crypto
ExecStart=/home/pi/algo-trade/crypto/.venv/bin/python dashboard.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
