# crypto/deploy/random_tick_buy.service
[Unit]
Description=Random Ticker 30-Min Strategy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/algo-trade/crypto
ExecStart=/home/pi/algo-trade/crypto/.venv/bin/python random_tick_buy.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
